# supervisor/puppetserver.conf
[program:puppetserver]
autostart=true
autorestart=true
command=/opt/puppetlabs/server/apps/puppetserver/bin/puppetserver foreground
process=%(program_name)s
